# tetsuo-lcd.service
[Unit]
Description=TETSUO LCD Price Display Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/tetsuo-display
Environment="PYTHONUNBUFFERED=1"
ExecStart=/usr/bin/python3 /home/pi/tetsuo-display/main_lcd.py --config config_lcd.yaml
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security settings
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target